# modules/admin/src/models/abryrath-backup.timer
[Unit]
Description=Site Backups

[Timer]
OnUnitActiveSec=60s
OnBootActiveSec=60s

[Install]
WantedBy=timers.target